# deploy/systemd/show-manager.service
[Unit]
Description=show-manager display controller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=devops
Group=devops
WorkingDirectory=/home/devops/workspace/projects/show-manager
Environment=HOME=/home/devops
EnvironmentFile=/home/devops/config/dev/show-manager.env
ExecStart=/usr/local/bin/node /home/devops/workspace/projects/show-manager/dist/server/src/main.js
Restart=always
RestartSec=3
TimeoutStopSec=20
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
